# services/golf-prizepicks-scraper.service
[Unit]
Description=Golf PrizePicks Scraper Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/GOLFAPP
ExecStart=/home/user/GOLFAPP/venv/bin/python -m services.prizepicks_scraper
Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/user/GOLFAPP/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=golf-pp-scraper

# Resource limits
MemoryMax=512M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
